# monitor.service
[Unit]
Description=Runs system monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor
