# systemmonitor.service
[Unit]
Description=Runs system monitor

[Service]
Type=oneshot
ExecStart=/usr/bin/systemmonitor collect --output database
